# dns-failover.service
[Unit]
Description=DNS Failover Monitor (Pi-hole to Authority)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dns-failover-check.sh
StandardOutput=journal
StandardError=journal
